# etc/supervisor-conf.d/supervisord-mysql.conf
[program:mysql]
command=start-mysql.sh
autorestart=true
autostart=true
